# Ultra96/packages/wilc3000/wpa_ap.service
[Unit]
Description=Wifi Access Point for U96-V2
Before=network-online.service isc-dhcp-server.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/share/wpa_ap/ap.sh

[Install]
WantedBy=basic.target
